# rc.local
#!/bin/bash

# do not change the name of the file in DOCKER, otherwise the container will not start

# Load AX25 Modules on the host in DOCKER
modprobe ax25
modprobe mkiss
modprobe netrom

# start openvpn (hamnet)
service openvpn start

# xinetd for uronode
service xinetd start

# custom startup
source /opt/startup/*.sh &


# Start completed, keep docker container alive (choose one of these)
tail -f /dev/null
sleep infinity
